# multi-node/vagrant/etcd2-spam.service
[Unit]
After=etcd2.service

[Service]
ExecStart=/opt/bin/etcd2-spam

[Install]
WantedBy=multi-user.target